# ProtonVPN/Linux/system/vpntime.service
[Unit]
Description=Force tyime-sync when using vpn

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpntime.sh
