# ProtonVPN/Linux/trigger/trigtimer.service
[Unit]
Description=Trigger vpntime.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start vpntime.service
